# roles/pulp_api/templates/pulpcore-api.socket.j2
[Unit]
Description=Pulp WSGI Server socket

[Socket]
ListenStream={{ pulp_api_bind | replace("unix:", "") }}

[Install]
WantedBy=sockets.target
