# roles/pulp_content/templates/pulpcore-content.socket.j2
[Unit]
Description=Pulp Content App socket

[Socket]
ListenStream={{ pulp_content_bind | replace("unix:", "") }}

[Install]
WantedBy=sockets.target
